# plists/dev.bash0c7.audio-transcription.swiftcap.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>            <string>dev.bash0c7.audio-transcription.swiftcap</string>
  <key>ProgramArguments</key> <array><string><%= swiftcap_bin %></string></array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SWIFTCAP_SPOOL</key>  <string><%= spool_dir %></string>
    <key>SWIFTCAP_LOCALE</key> <string>ja-JP</string>
  </dict>
  <key>RunAtLoad</key>     <true/>
  <key>KeepAlive</key>     <true/>
  <key>StandardOutPath</key> <string><%= log_dir %>/swiftcap.log</string>
  <key>StandardErrorPath</key><string><%= log_dir %>/swiftcap.err</string>
</dict>
</plist>
